# config/systemd/taskforce-development-debug.service
[Unit]
Description=taskforce local web UI (development, debug)
After=network.target

[Service]
Type=simple
Environment=TASKFORCE_ENV=development
EnvironmentFile=-%h/.config/taskforce/taskforce.env
ExecStart=%h/work/myrepo/taskforce/target/debug/taskforce serve
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
